# scripts/olms-go.cron
#! /bin/bash

/var/www/olms-go/olms-go backup